# gps-stream.service
[Unit]
Description=Enhanced Garmin GPS Broadcaster with Auto-Detection
After=network-online.target
Wants=network-online.target
Before=gpsd.service
Conflicts=gpsd.service gpsd.socket

[Service]
Type=notify
ExecStartPre=/home/cuas/gps_startup_cleanup.sh
ExecStart=/usr/bin/python3 /home/cuas/garminReader.py
WorkingDirectory=/home/cuas
User=root
Group=root
EnvironmentFile=-/etc/default/gps-stream
Environment=PYTHONUNBUFFERED=1





# Enhanced restart behavior
Restart=always
RestartSec=15
StartLimitBurst=5
StartLimitIntervalSec=300

# Watchdog configuration
Type=notify
WatchdogSec=60
NotifyAccess=main

# Resource limits for Pi 2
MemoryMax=256M

# Process management
KillMode=control-group
KillSignal=SIGTERM
TimeoutStartSec=120
TimeoutStopSec=30

# Enhanced logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gps-stream

[Install]
WantedBy=multi-user.target
