# gps-system-watchdog.service
[Unit]
Description=GPS System Ultimate Watchdog - Last Line of Defense
After=gps-stream.service
Wants=gps-stream.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/cuas/system_watchdog.py
Restart=always
RestartSec=30
User=root
Group=root
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
